# salt/systemd/systemd.service.jinja
[Unit]
Description=todo
After=syslog.target network.target

[Service]
Type=simple
RemainAfterExit=no
User={{ user }}
Group={{ group}}
ExecStart={{ execstart }} 

[Install]
WantedBy=multi-user.target